# x-ui/xui-watcher.service
[Unit]
Description=XUI DB Watcher and Location Sync
After=network.target

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/xui-watcher.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
